# conf/supervisor.conf
[program:personal_site]
command=/home/ubuntu/personal_site/run.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/personal_site.err.log
stdout_logfile=/var/log/personal_site.out.log
